# dist/nectar-server.service
[Unit]
Description=Nectar-Server Daemon
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/nectar-server

[Install]
WantedBy=multi-user.target
